# content/etc/s6-overlay/s6-rc.d/powerdns/run
!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Starting powerdns"
exec s6-setuidgid razor /usr/sbin/pdns_server 